# assets/etc/supervisor/conf.d/tomcat.ini
[program:tomcat]
environment=CATALINA_HOME=%(ENV_CATALINA_HOME)s
command=/bin/sh -c '${CATALINA_HOME}/bin/catalina.sh run'

autostart=true
autorestart=true
stopasgroup=true
startasgroup=true
